# security-knowledge/ops/security-knowledge-worker.service
[Unit]
Description=Security Knowledge ARQ worker
After=network.target postgresql.service redis.service security-knowledge.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=z
Group=z
WorkingDirectory=/home/z/z-zoidberg/security-knowledge
EnvironmentFile=/home/z/z-zoidberg/security-knowledge/.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/home/z/z-zoidberg/security-knowledge/.venv/bin/arq app.worker.WorkerSettings
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
